# cron/crontab
* * * * * /bin/bash /home/pi/rpi-photo-frame/scripts/update.sh >> /var/log/updater.log
@reboot /usr/bin/python3 /home/pi/rpi-photo-frame/server/manage.py >> /var/log/photo-frame.log
