# cron/crontab
* * * * * /bin/bash /home/pi/rpi-photo-frame/scripts/update.sh >> /dev/null 2>&1
@reboot /bin/bash /home/pi/rpi-photo-frame/scripts/kiosk.sh >> /dev/null 2>&1
