# cron/crontab
* * * * * /bin/sh /home/pi/rpi-photo-frame/src/scripts/update.sh >> /dev/null 2>&1