# cron/crontab
* * * * * /bin/sh /home/pi/rpi-photo-frame/scripts/update.sh >> /dev/null 2>&1
@reboot export PHOTO_FRAME_ENV=dev && /usr/bin/python3.5 /home/pi/rpi-photo-frame/server/manage.py >> /dev/null 2>&1
