# cron/crontab
* * * * * /bin/bash /home/pi/rpi-photo-frame/scripts/update.sh >> /var/log/updater.log
