# scripts/rc.local
#!/bin/sh -e

/bin/bash /home/pi/rpi-photo-frame/scripts/kiosk.sh

dmesg --console-off
clear

exit 0
