# src/conf/datenwissenschaften-splash.service
[Unit]
Description=Starts Datenwissenschaften Splash Screen
DefaultDependencies=true

[Service]
Type=simple
ExecStart=/bin/sh /home/pi/rpi-photo-frame/src/scripts/bootstrap.sh >> /var/log/bootstrap.log
WorkingDirectory=/home/pi/

[Install]
WantedBy=multi-user.target