# src/conf/rc.local
#!/bin/sh -e

# Turn all console prints off
dmesg --console-off &

# Show splash screen
/usr/bin/fbi -T 1 -noverbose -a -t 30 --once /home/pi/rpi-photo-frame/doc/splash.png &

# Clear console
clear

exit 0