# src/conf/rc.local
#!/bin/sh -e

# Turn backlight off
# /usr/local/bin/rpi-backlight --off &

# Turn all console prints off
dmesg --console-off &

# Show splash screen
/usr/bin/fbi -T 3 -noverbose -a -t 30 --once /home/pi/rpi-photo-frame/doc/splash.png &

# Clear console
clear

exit 0